# supervisor-autofs.conf
[program:autofs]
command=/usr/sbin/automount -t 0 -f /etc/auto.master
numprocs=1
autostart=false
autorestart=false
startretries=3
