# supervisor-autofs.conf
[program:autofs]
command=/usr/sbin/automount -v -t 0 -f /etc/auto.master
numprocs=1
autostart=true
autorestart=true
startretries=3
