# supervisor-autofs.conf
[program:autofs]
command=BROWSE_MODE="no" /usr/sbin/automount -t 0 -f /etc/auto.master
numprocs=1
autostart=true
autorestart=true
startretries=3
user=root
